# cix_proprietary/cix_proprietary-debs/cix-isp-umd/etc/systemd/system/load-isp-modules.service
[Unit]
Description=Load isp kernel modules

[Service]
Type=oneshot
ExecStart=/usr/bin/load-isp-modules.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
