# cix_proprietary/cix_proprietary-debs/cix-gpu-umd/etc/systemd/system/load-gpu-modules.service
[Unit]
Description=Load gpu modules

[Service]
Type=oneshot
ExecStart=/usr/bin/load-gpu-modules.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
